# nix/packages/unifi-os-server-image/rootfs/etc/systemd/system/mongodb.service.d/mkdirs.conf
[Service]
ExecStartPre=+/bin/mkdir -p /var/log/mongodb /var/lib/mongodb
ExecStartPre=+/bin/chown mongodb:mongodb /var/log/mongodb /var/lib/mongodb
