# nix/packages/unifi-os-server-image/rootfs/etc/systemd/system/unifi-core.service.d/mkdirs.conf
[Service]
ExecStartPre=+/bin/mkdir -p /data/unifi-core/config/http
